# scripts/systemd/catscan-contracts-check.timer
[Unit]
Description=Daily Cat-Scan contract validation

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
