# scripts/systemd/catscan-home-refresh.timer
[Unit]
Description=Nightly Cat-Scan Home cache refresh

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
