# scripts/systemd/catscan-creative-cache-refresh.timer
[Unit]
Description=Nightly Cat-Scan creative cache refresh

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true

[Install]
WantedBy=timers.target

